# crontab
*/2 * * * * python /home/pi/infraspecific/log_video.py
*/2 * * * * python /home/pi/infraspecific/log_audio.py
0 3 * * * python /home/pi/infraspecific/move_logs.py